# camera/netwatch.service
[Unit]
Description=Network self-heal watchdog (bird-cam)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/netwatch.sh
